# supervisord/conf.d/artifactory.conf
[program:artifactory]
command=/bin/bash -c "/etc/init.d/artifactory start"
redirect_stderr=true