# supervisord/conf.d/jenkins.conf
[program:jenkins]
command=/bin/bash -c "/etc/init.d/jenkins start"
redirect_stderr=true