# supervisord/conf.d/sonar.conf
[program:sonar]
command=/bin/bash -c "/etc/init.d/sonar start"
redirect_stderr=true