# supervisord/conf.d/tomcat.conf
[program:simple-webapp]
command=/bin/bash -c "/etc/init.d/tomcat start"
redirect_stderr=true
